# daemon/pkg/post-install.sh
#!/bin/bash
systemctl enable aerospike-local-daemon.service
