# aerospike-connection-manager/pkg/build/usr/lib/systemd/system/aerospike-connection-manager.service
[Unit]
Description=Aerospike Local Daemon Service
Documentation=https://https://github.com/aerospike/php-client/asld
Wants=network.target
After=network-online.target

[Service]
ExecStart=/usr/bin/aerospike-connection-manager --config /etc/aerospike-connection-manager/asld.toml

[Install]
WantedBy=multi-user.target
