# aerospike-connection-manager/pkg/post-install.sh
#!/bin/bash
systemctl enable aerospike-connection-manager.service
